# deploy/systemd/chum.service
[Unit]
Description=CHUM Agent Orchestrator
After=openclaw-gateway.service

[Service]
Type=simple
ExecStart=%h/projects/chum/chum --config %h/projects/chum/chum.toml
Restart=always
RestartSec=10
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin

[Install]
WantedBy=default.target
