# deploy/systemd/chum.service
[Unit]
Description=CHUM Agent Orchestrator
After=openclaw-gateway.service

[Service]
Type=simple
ExecStartPre=/home/ubuntu/projects/cortex/chum stop
ExecStart=/home/ubuntu/projects/cortex/chum --config /home/ubuntu/projects/cortex/chum.toml
Restart=always
RestartSec=10
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin

[Install]
WantedBy=default.target
